# openbazaar_template.conf
# Startup script for OpenBazaar Server that will auto-restart if crashed
# Copy this script to /etc/init folder
# "chmod 644 openbazaar.conf"
# Change "chdir", "setuid" and "setgid" according to your system
# Usage: "sudo service openbazaar start"

description "OpenBazaar Server"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

chdir /root/OpenBazaar-Server 

setuid root
setgid root

respawn
respawn limit 10 5
kill timeout 20

console log

#exec_template
